# lib/phantomjs-highchart-pin/phantomjs-highchart-pin.conf
description 'phantomjs-highchart'
start on runlevel [2345]
stop on runlevel [06]
respawn
# in case of failure retry 3 times to respawn with 5 sec interval
respawn limit 3 5

# set some variables
env SERVER=/usr/local/share/phantomjs-highchart-pin/highcharts-convert.js
env HOST=127.0.0.1
env PORT=3003

# enable logger output and write to syslog with tag phantomjs
# for more info see http://upstart.ubuntu.com/cookbook/#id152
# and/or type logger --help
console output
# execute the command and log it
exec phantomjs $SERVER -host $HOST -port $PORT | logger -s -t "phantomjs: "
